# ubuntu/recapture-docs.conf
# Ubuntu upstart file at /etc/init/recapturedocs.conf

start on runlevel [2345]
stop on runlevel [06]

script
  exec /home/jaraco/envs/recapturedocs/bin/recapture-docs serve -C prod
end script
